# s6/setup-dirs/up
#!/command/with-contenv sh
s6-mkdir -p /data/pb_data