# s6/api-server/run
#!/command/with-contenv sh
exec /app/server/serve